# service/agent/run
#!/bin/sh
: ${GOOSE="/opt/goose"}
cd ${GOOSE}/goose
chmod g-w ${GOOSE}/goose ${GOOSE}/goose/config
touch ${GOOSE}/goose/config/authorized_keys
exec chpst -e ${GOOSE}/service/agent/env ${GOOSE}/bin/shell vendor/projects/offline/bin/setup-ssh-agent ${GOOSE}/goose/config/ssh_agent ${GOOSE}/goose/config/authorized_keys
