# service/sshd/run
#!/bin/sh
cd ${GOOSE}
exec chpst -e ${GOOSE}/service/sshd/env ${GOOSE}/vendor/projects/offline/bin/setup-sshd ${GOOSE}/config/sshd.port ${GOOSE}/config/authorized_keys
