# service/proxy/run
#!/bin/sh
cd ${GOOSE}
exec chpst -e ${GOOSE}/service/proxy/env ${GOOSE}/vendor/projects/offline/bin/setup-squid ${GOOSE}/config/proxy.port
