# service/proxy/run
#!/bin/sh
: ${GOOSE="/opt/goose"}
cd ${GOOSE}/goose
exec chpst -e ${GOOSE}/service/proxy/env ${GOOSE}/bin/shell vendor/projects/offline/bin/setup-squid ${GOOSE}/goose/config/proxy.port
